# src/webhook-server.service
[Unit]
Description=Convert Chat Webhook Server
After=network.target
Documentation=https://github.com/aibro-ctrl/Convert

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/Convert/src
Environment="NODE_ENV=production"
Environment="WEBHOOK_PORT=9000"
Environment="WEBHOOK_SECRET=your-secret-here-change-me"
ExecStart=/usr/bin/node /var/www/Convert/src/webhook-server.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/convert-webhook.log
StandardError=append:/var/log/convert-webhook-error.log

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/Convert/src /tmp /var/log

[Install]
WantedBy=multi-user.target
